# Petalinux 2019.2/WithY3_UART1-to-Data2_Ubuntu16/rootfs/var/lib/dpkg/info/acpid.conffiles
/etc/acpi/events/powerbtn
/etc/acpi/powerbtn.sh
/etc/default/acpid
/etc/init.d/acpid
/etc/init/acpid.conf
