# Petalinux 2019.2/WithY3_UART1-to-Data2_Ubuntu16/rootfs/var/lib/dpkg/info/connman.conffiles
/etc/apport/crashdb.conf.d/connman-crashdb.conf
/etc/dbus-1/system.d/connman.conf
/etc/init.d/connman
/etc/init/connman.conf
