# Petalinux 2019.2/WithY3_UART1-to-Data2_Ubuntu16/rootfs/var/lib/dpkg/info/connman.list
/.
/etc
/etc/init.d
/etc/init.d/connman
/etc/apport
/etc/apport/crashdb.conf.d
/etc/apport/crashdb.conf.d/connman-crashdb.conf
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/connman.conf
/etc/init
/etc/init/connman.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/connman.service
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/connmand.8.gz
/usr/share/man/man8/connman.8.gz
/usr/share/man/man5
/usr/share/man/man5/connman.conf.5.gz
/usr/share/man/man1
/usr/share/man/man1/connmanctl.1.gz
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_connman.py
/usr/share/doc
/usr/share/doc/connman
/usr/share/doc/connman/changelog.Debian.gz
/usr/share/doc/connman/copyright
/usr/share/doc/connman/README.Debian
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/net.connman.policy
/usr/sbin
/usr/sbin/connmand
/usr/sbin/connmanctl
