# Petalinux 2019.2/WithY3_UART1-to-Data2_Ubuntu16/rootfs/var/lib/dpkg/info/cpufrequtils.list
/.
/usr
/usr/share
/usr/share/locale
/usr/share/locale/pt
/usr/share/locale/pt/LC_MESSAGES
/usr/share/locale/pt/LC_MESSAGES/cpufrequtils.mo
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/cpufrequtils.mo
/usr/share/locale/it
/usr/share/locale/it/LC_MESSAGES
/usr/share/locale/it/LC_MESSAGES/cpufrequtils.mo
/usr/share/locale/fr
/usr/share/locale/fr/LC_MESSAGES
/usr/share/locale/fr/LC_MESSAGES/cpufrequtils.mo
/usr/share/locale/ca
/usr/share/locale/ca/LC_MESSAGES
/usr/share/locale/ca/LC_MESSAGES/cpufrequtils.mo
/usr/share/locale/cs
/usr/share/locale/cs/LC_MESSAGES
/usr/share/locale/cs/LC_MESSAGES/cpufrequtils.mo
/usr/share/doc
/usr/share/doc/cpufrequtils
/usr/share/doc/cpufrequtils/README.Debian
/usr/share/doc/cpufrequtils/copyright
/usr/share/doc/cpufrequtils/examples
/usr/share/doc/cpufrequtils/examples/cpufrequtils.loadcpufreq.sample
/usr/share/doc/cpufrequtils/examples/cpufrequtils.sample
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/cpufrequtils
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/cpufreq-set.1.gz
/usr/share/man/man1/cpufreq-info.1.gz
/usr/share/man/man1/cpufreq-aperf.1.gz
/usr/bin
/usr/bin/cpufreq-set
/usr/bin/cpufreq-info
/usr/bin/cpufreq-aperf
/etc
/etc/init.d
/etc/init.d/loadcpufreq
/etc/init.d/cpufrequtils
/usr/share/doc/cpufrequtils/changelog.Debian.gz
