# Petalinux 2019.2/WithY3_UART1-to-Data2_Ubuntu16/rootfs/var/lib/dpkg/info/hostapd.list
/.
/usr
/usr/sbin
/usr/sbin/hostapd_cli
/usr/sbin/hostapd
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/hostapd
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/hostapd.8.gz
/usr/share/man/man1
/usr/share/man/man1/hostapd_cli.1.gz
/usr/share/doc
/usr/share/doc/hostapd
/usr/share/doc/hostapd/examples
/usr/share/doc/hostapd/examples/hostapd.deny
/usr/share/doc/hostapd/examples/hostapd.radius_clients
/usr/share/doc/hostapd/examples/hostapd.eap_user.gz
/usr/share/doc/hostapd/examples/hostapd.accept
/usr/share/doc/hostapd/examples/hostapd.conf.gz
/usr/share/doc/hostapd/examples/hostapd.wpa_psk
/usr/share/doc/hostapd/copyright
/usr/share/doc/hostapd/changelog.Debian.gz
/usr/share/doc/hostapd/README.Debian
/usr/share/doc/hostapd/NEWS.Debian.gz
/etc
/etc/init.d
/etc/init.d/hostapd
/etc/network
/etc/network/if-pre-up.d
/etc/network/if-post-down.d
/etc/default
/etc/default/hostapd
/etc/hostapd
/etc/hostapd/ifupdown.sh
/etc/network/if-pre-up.d/hostapd
/etc/network/if-post-down.d/hostapd
