# Petalinux 2019.2/WithY3_UART1-to-Data2_Ubuntu16/rootfs/var/lib/dpkg/info/hostapd.conffiles
/etc/default/hostapd
/etc/hostapd/ifupdown.sh
/etc/init.d/hostapd
