# Petalinux 2019.2/WithY3_UART1-to-Data2_Ubuntu16/rootfs/var/lib/dpkg/info/ofono.list
/.
/usr
/usr/sbin
/usr/sbin/ofonod
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/ofono
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/ofonod.8.gz
/usr/share/doc
/usr/share/doc/ofono
/usr/share/doc/ofono/copyright
/usr/share/doc/ofono/AUTHORS.gz
/usr/share/doc/ofono/changelog.Debian.gz
/var
/var/lib
/var/lib/ofono
/etc
/etc/ofono
/etc/ofono/phonesim.conf
/etc/init
/etc/init/ofono.conf
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/ofono.conf
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/97-ofono-speedup.rules
/lib/udev/rules.d/97-ofono.rules
/lib/systemd
/lib/systemd/system
/lib/systemd/system/ofono.service
