# Petalinux 2019.2/WithY3_UART1-to-Data2_Ubuntu16/rootfs/var/lib/dpkg/info/udhcpd.list
/.
/usr
/usr/sbin
/usr/bin
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man1
/usr/share/man/man5
/usr/share/man/man5/udhcpd.conf.5.gz
/usr/share/doc
/usr/share/doc/udhcpd
/usr/share/doc/udhcpd/copyright
/etc
/etc/init.d
/etc/init.d/udhcpd
/etc/default
/etc/default/udhcpd
/etc/udhcpd.conf
/usr/sbin/udhcpd
/usr/bin/dumpleases
/usr/share/man/man8/udhcpd.8.gz
/usr/share/man/man1/dumpleases.1.gz
/usr/share/doc/udhcpd/changelog.Debian.gz
