# Petalinux 2019.2/WithY3_UART1-to-Data2_Ubuntu16/rootfs/var/lib/dpkg/info/udhcpd.conffiles
/etc/default/udhcpd
/etc/init.d/udhcpd
/etc/udhcpd.conf
