# Petalinux 2019.2/WithY3_UART1-to-Data2_Ubuntu16/rootfs/lib/systemd/system/apache2.service.d/apache2-systemd.conf
[Service]
Type=forking
RemainAfterExit=no
